# etc/crontab/cron
# m h  dom mon dow   command
#* * * * * root /usr/local/bin/php 
* * * * * root /usr/local/bin/php /var/www/html/goldpillar/release/bin/console app:generate:monthlyfee >> /var/log/goldpillar-app.generate.monthlyfee.log